# rootfs/etc/services.d/openlitespeed/run
#!/usr/bin/with-contenv bash

if [ -f /xshok-init.sh ] ; then
  echo "WOOF I" >> /tmp/meow
  while ! bash /xshok-init.sh ; do
    echo "Waiting for init"
    sleep 5s
  done
fi
