# rootfs/etc/services.d/cron/run
#!/bin/execlineb -P

exec cron -f -L 15
