# rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv bash
################################################################################
# This is property of eXtremeSHOK.com
# You are free to use, modify and distribute, however you may not remove this notice.
# Copyright (c) Adrian Jon Kriel :: admin@extremeshok.com
################################################################################

# Runs cron every 1 minute
# Generates a new vhost cron every 15mins

while ! /usr/local/lsws/bin/lswsctrl status | grep -q "litespeed is running with PID" ; do
	echo "Waiting for OpenLiteSpeed to start"
  sleep 10s
done

runcount=0
while true ; do
  if [ "$runcount" -eq 0 ] || [ "$runcount" -eq 15 ] ; then
    if [ -f /xshok-vhost-generate-cron.sh ] ; then
      echo "Launching xshok-vhost-generate-cron.sh"
      bash /xshok-vhost-generate-cron.sh
    fi
    runcount=1
  else
    (( runcount++ ))
  fi
  cron -f
  sleep 60s
done
