# rootfs/etc/cron.hourly/generate-vhost-cron
#!/usr/bin/env bash

echo "cron.hourly" >> /tmp/testing

bash /xshok-generate-vhost-cron.sh
