# rootfs/etc/cron.hourly/generate-vhost-cron
#!/usr/bin/env bash

bash /xshok-generate-vhost-cron.sh
