# rootfs/etc/services.d/inotify-certs/run
#!/usr/bin/with-contenv bash
################################################################################
# This is property of eXtremeSHOK.com
# You are free to use, modify and distribute, however you may not remove this notice.
# Copyright (c) Adrian Jon Kriel :: admin@extremeshok.com
################################################################################

# restart openlitespeed when changes to the vhost/domain.com/cert dirs are detected

# this is a work around to support files in new dirs (when a new website is added, it was not being watched)
# double negative, so the exclude becomes an only include, only will include *.pem files in the non excluded directories (certs)

while ! /usr/local/lsws/bin/lswsctrl status | grep -q "litespeed is running with PID" ; do
	echo "Waiting for OpenLiteSpeed to start"
  sleep 10s
done

while true ; do
	inotifywait -r --exclude "\.[^p].[^e].[^m]$" --exclude ".*(\.git|\.private|private|html|public_html|www|db|dbinfo|log|logs|sql|backup|backups|conf|config|configs)/" /var/www/vhosts/ -e modify -e create -e delete -e move
	echo "sleeping for 60s"
	sleep 60s
	/usr/local/lsws/bin/lswsctrl restart &
done
